# server/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

cd /usr/lib/python2.4/site-packages/xmppd/; python xmppd.py -d > xmppd.output &
sleep 10
cd /var/www/html/myweb/; python chatServer.py &
